# s6-overlay/services.d/php/finish
#!/bin/sh

echo "PHP service terminated with exit code $1"
exit "$1"
