# s6-overlay/services.d/php/run
#!/bin/sh
# Note: Changed from bash to sh for better compatibility

# This script determines how to start the PHP service based on the image type

# Detect PHP type by checking for binaries (more reliable than SAPI detection)
if [ -z "${PHP_TYPE}" ]; then
    if command -v php-fpm >/dev/null 2>&1; then
        PHP_TYPE="fpm"
    elif command -v apache2-foreground >/dev/null 2>&1 || command -v apache2 >/dev/null 2>&1; then
        PHP_TYPE="apache"
    else
        PHP_TYPE="cli"
    fi
fi

echo "Starting PHP service in ${PHP_TYPE} mode..."

case "$PHP_TYPE" in
  "cli")
    # For CLI, we just execute a long-running command or wait
    # This is typically handled by the CMD instruction and not S6
    echo "Running in CLI mode - service management handled by CMD"
    sleep infinity
    ;;
    
  "fpm")
    # Start PHP-FPM in foreground mode
    # Note: php-fpm master runs as root and drops privileges via pool config
    echo "Starting PHP-FPM..."
    exec php-fpm -F
    ;;
    
  "apache")
    # Start Apache in foreground mode
    echo "Starting Apache..."
    exec apache2-foreground
    ;;
    
  *)
    echo "Unknown PHP type: ${PHP_TYPE}, defaulting to CLI mode"
    sleep infinity
    ;;
esac
