# s6-overlay/services.d/php/run
#!/bin/sh
# Note: Changed from bash to sh for better compatibility

# This script determines how to start the PHP service based on the image type

# Detect PHP type from container environment or image tag
PHP_TYPE=${PHP_TYPE:-$(php -r 'echo (php_sapi_name() === "cli" ? "cli" : (php_sapi_name() === "fpm" ? "fpm" : "apache"));')}

echo "Starting PHP service in ${PHP_TYPE} mode..."

case "$PHP_TYPE" in
  "cli")
    # For CLI, we just execute a long-running command or wait
    # This is typically handled by the CMD instruction and not S6
    echo "Running in CLI mode - service management handled by CMD"
    sleep infinity
    ;;
    
  "fpm")
    # Start PHP-FPM with appropriate user
    echo "Starting PHP-FPM..."
    exec s6-setuidgid appuser php-fpm
    ;;
    
  "apache")
    # Start Apache with appropriate user
    echo "Starting Apache..."
    exec apache2-foreground
    ;;
    
  *)
    echo "Unknown PHP type: ${PHP_TYPE}, defaulting to CLI mode"
    sleep infinity
    ;;
esac
